# config/.config/systemd/user/sync-google-drive-rclone.service
[Unit]
Description=Sync google drive with rclone every 5 minutes.

[Service]
ExecStart=%h/.scripts/cloudsync
Type=oneshot
